# scripts/supervisor/cm.conf
[program:cloudera-scm-server-db]
command=/bin/bash -c "/etc/init.d/cloudera-scm-server-db start"

[program:cloudera-scm-server]
command=/bin/bash -c "/etc/init.d/cloudera-scm-server start"
